# memory-swap/claude-spare-reaper.service
[Unit]
Description=Reap orphaned Claude Code spare processes (daemon died, nothing running inside)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/claude-spare-reaper
