# wifi/etc/systemd/system/wifi-recover.service
[Unit]
Description=Recover Intel AX201 WiFi after firmware death
After=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/wifi-recover
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
